# infra/prod/deploy/pulse-secret-sync.timer
[Unit]
Description=Check PULSE Secrets Manager values every five minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
RandomizedDelaySec=30s
Persistent=true

[Install]
WantedBy=timers.target
