# infra/prod/deploy/pulse-secret-sync.service
[Unit]
Description=PULSE Secrets Manager runtime sync
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/pulse-refresh-secrets
